# {{ cookiecutter.project_directory }}/supervisor.conf
[program:{{ cookiecutter.project_slug }}]
command=/var/apps/{{ cookiecutter.project_directory }}/gunicorn_start
stopasgroup=true
killasgroup=true
